# task.xml
<?xml version='1.0' encoding='UTF-16'?>
<ns0:Task xmlns:ns0="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">

  <ns0:RegistrationInfo>

    <ns0:Date>2021-10-19T11:21:41</ns0:Date>

    <ns0:Author>DESKTOP-ALQA6CL\sohum</ns0:Author>

    <ns0:URI>\template</ns0:URI>

  </ns0:RegistrationInfo>

  <ns0:Principals>

    <ns0:Principal id="Author">

      <ns0:UserId>S-1-5-21-61245436-1204012862-1097636640-1004</ns0:UserId>

      <ns0:LogonType>Password</ns0:LogonType>

      <ns0:RunLevel>HighestAvailable</ns0:RunLevel>

    </ns0:Principal>

  </ns0:Principals>

  <ns0:Settings>

    <ns0:DisallowStartIfOnBatteries>false</ns0:DisallowStartIfOnBatteries>

    <ns0:StopIfGoingOnBatteries>true</ns0:StopIfGoingOnBatteries>

    <ns0:MultipleInstancesPolicy>IgnoreNew</ns0:MultipleInstancesPolicy>

    <ns0:StartWhenAvailable>true</ns0:StartWhenAvailable>

    <ns0:WakeToRun>true</ns0:WakeToRun>

    <ns0:IdleSettings>

      <ns0:StopOnIdleEnd>true</ns0:StopOnIdleEnd>

      <ns0:RestartOnIdle>false</ns0:RestartOnIdle>

    </ns0:IdleSettings>

  </ns0:Settings>

  <ns0:Triggers>

    <ns0:TimeTrigger>

      <ns0:StartBoundary>2022-03-21T6:59:50</ns0:StartBoundary>

    </ns0:TimeTrigger>

  </ns0:Triggers>

  <ns0:Actions Context="Author">

    <ns0:Exec>

      <ns0:Command>"C:\Users\sohum\OneDrive\Documents\Python Scripts\Selenium\tennis\runTennisAutomation.bat"</ns0:Command>

      <ns0:Arguments>2022-03-24 13 14</ns0:Arguments>

      <ns0:WorkingDirectory>C:\Users\sohum\OneDrive\Documents\Python Scripts\Selenium</ns0:WorkingDirectory>

    </ns0:Exec>

  </ns0:Actions>

</ns0:Task>